# overlay/mesa-freedreno-nobin/etc/profile.d/mesa-freedreno-nobin.sh
# Hacky workaround to make GTK4/GNOME not artifact like crazy on broken A5XX freedreno(/kernel?) drivers
export FD_MESA_DEBUG=nobin
